# PUSH-deployment-fixed/pkg-build/payload/Library/LaunchDaemons/com.push.autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.push.autoupdate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Management/PUSH/push-cli</string>
        <string>auto-check</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Library/Management/PUSH/logs/push-cli.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Management/PUSH/logs/push-cli.log</string>
    <key>KeepAlive</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
